# roles/app/templates/supervisor.conf
{% for worker in workers %}
[program:mobili_worker_{{ worker }}]
command = python3 manage.py runworker ;
directory = {{ dir }}/backend/src ;
user = {{ user }} ;
stdout_logfile = {{ dir }}/logs/workers.log ;
redirect_stderr = true ;
autostart = true ;
autorestart = true ;
killasgroup = true ;

{% endfor %}

[program:mobili_daphne]
command = daphne mobili.asgi:channel_layer --unix-socket {{ dir }}/daphne.socket ;
directory = {{ dir }}/backend/src ;
user = {{ user }} ;
stdout_logfile = {{ dir }}/logs/daphne.log ;
redirect_stderr = true ;
autostart = true ;
autorestart = true ;
killasgroup = true ;
